# systemd/smartserver.service
# /etc/systemd/system/smartserver.service
[Unit]
Description=Smart Home - API Server
Wants=network-online.target
After=network-online.target mosquitto.service

[Service]
WorkingDirectory=/home/adams/zigbeemcp
User=adams
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/adams/zigbeemcp/.env

# Use venv python to run the API server
ExecStart=/home/adams/zigbeemcp/.venv/bin/python -m smarthome.api.server

Restart=always
RestartSec=2
StartLimitIntervalSec=30
StartLimitBurst=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
